# supervisord.conf
[supervisord]
nodaemon=true

[program:nodejs]
directory=/var/www
command=npm start
autorestart = true
stderr_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:avgd]
directory=/var/www
command=/opt/avg/av/bin/avgd
autorestart = true
environment=AVG_CTL='avgctl',AVG_CTL_PATH='/opt/avg/av/bin/avgctl',AVG_DAEMON='avgd',AVG_DESC='AVG Anti-Virus Daemon',AVG_DIAG='avgdiag',AVG_LOCK_FILE='/var/lock/avgd',AVG_LOCK_FILE_DIR='/var/lock',AVG_LOCK_FILE_NAME='avgd',AVG_PID_FILE='/opt/avg/av/var/run/avgd.pid',AVG_PROG='avgd',AVG_RUN_DIR='/opt/avg/av/var/run',AVG_SCRIPTS_DIR='/opt/avg/av/etc/init.d',AVG_WAIT_FOR_PID_TIME='15',DST_TREE='/opt/avg/av'